# config/systemd/okimi-auth.service
[Unit]
Description=Okimi Auth Service
After=network.target postgresql.service redis.service keycloak.service
Wants=postgresql.service redis.service keycloak.service

[Service]
Type=simple
User=okimi
Group=okimi
WorkingDirectory=/opt/okimi/services/auth
Environment="JAVA_OPTS=-Xmx512m -Xms256m"
Environment="SPRING_PROFILES_ACTIVE=production"
ExecStart=/usr/bin/java $JAVA_OPTS -jar /opt/okimi/services/auth/auth-service.jar
Restart=on-failure
RestartSec=10s

# Resource limits
MemoryMax=1G
CPUQuota=200%

# Security
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=okimi-auth

[Install]
WantedBy=multi-user.target
