# config/systemd/okimi-gateway.service
[Unit]
Description=Okimi API Gateway
After=network.target okimi-auth.service okimi-user.service
Wants=okimi-auth.service okimi-user.service

[Service]
Type=simple
User=okimi
Group=okimi
WorkingDirectory=/opt/okimi/services/gateway
Environment="JAVA_OPTS=-Xmx512m -Xms256m"
Environment="SPRING_PROFILES_ACTIVE=production"
ExecStart=/usr/bin/java $JAVA_OPTS -jar /opt/okimi/services/gateway/gateway-service.jar
Restart=on-failure
RestartSec=10s

# Resource limits
MemoryMax=1G
CPUQuota=200%

# Security
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=okimi-gateway

[Install]
WantedBy=multi-user.target
